# aview.service.example
[Unit]
Description=A-View Document Viewer Service
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=kdy987
Group=kdy987
WorkingDirectory=/opt/aview
ExecStart=/opt/aview/aview
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# 환경변수 설정
Environment=AVIEW_MODE=vx50
Environment=PORT=8003
Environment=HOST=0.0.0.0

# 보안 설정
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/aview/logs /opt/aview/cache /opt/aview/db /opt/aview/files

[Install]
WantedBy=multi-user.target